# docs/testdeploy_test.ini
[program:testdeploy]
directory=/home/web/testdeploy/
command=/home/web/testdeploy/testdeploy -config=/home/web/testdeploy/test.toml
stdout_logfile=/data/log/testdeploy/testdeploy.log
stdout_logfile_backups=50
redirect_stderr=true
autostart=true
autorestart=true